# resources/systemd/backup-nextcloud.timer
[Unit]
Description=Run regular backups of nextcloud

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
